# claude-session-starter@.service
[Unit]
Description=Claude Session Starter - Smart AI Session Management
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/MrMohebi/claude-session-starter

[Service]
Type=simple
User=%i
Group=%i

# Load installation path from environment file
# Create /etc/default/claude-session-starter-%i with:
# INSTALL_PATH=/path/to/claude-session-starter
EnvironmentFile=-/etc/default/claude-session-starter-%i

# Main service command (uses INSTALL_PATH from environment file, defaults to /home/%i/claude-session-starter)
ExecStart=/bin/sh -c 'exec ${INSTALL_PATH:-/home/%i/claude-session-starter}/claude-session-starter.sh --daemon'

# Restart policy
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/.config/claude-session-starter /home/%i/.claude

# Environment
Environment="PATH=/home/%i/.local/bin:/usr/local/bin:/usr/bin:/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-session-starter

[Install]
WantedBy=multi-user.target
